# server/systemd/strataslims.service
[Unit]
Description=StrataSlims Discord Bot
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=korokchatterjee
Group=korokchatterjee
WorkingDirectory=/home/korokchatterjee/StrataSlims

# Use the project venv’s Python and run the bot entry point

ExecStart=/home/korokchatterjee/StrataSlims/.venv/bin/python /home/korokchatterjee/StrataSlims/main.py
# ExecStart=/bin/bash -lc '/home/korokchatterjee/StrataSlims/server/run_with_update.sh'
ExecStopPost=/usr/bin/env python3 /home/korokchatterjee/StrataSlims/server/notify_preemption.py

# Ensure stdout/err go to journal

StandardOutput=journal
StandardError=journal

# Auto-restart on crash or disconnect

Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=20

# Optional: export .env to the process environment (your code also loads it via python-dotenv)

EnvironmentFile=/home/korokchatterjee/StrataSlims/.env
Environment=PYTHONUNBUFFERED=1
Environment=LC_ALL=C.UTF-8
Environment=LANG=C.UTF-8

[Install]
WantedBy=multi-user.target