# strataslims.conf
[program:strataslims]
; Entry point performs git pull then starts the bot. Adjust paths to match your VM.
command=/bin/bash -lc '/home/korokchatterjee/StrataSlims/run_with_update.sh'
directory=/home/korokchatterjee/StrataSlims
user=korokchatterjee

; Start and restart behavior
autostart=true
autorestart=true
startsecs=5
startretries=10
stopsignal=TERM
stopasgroup=true
killasgroup=true

; Environment variables for the bot. Prefer .env file but you can inline critical ones.
; Example (uncomment and replace values):
; environment=DISCORD_BOT_TOKEN="${DISCORD_BOT_TOKEN}",TEST_GUILD_ID="${TEST_GUILD_ID}",GREENLIST_USER_IDS="${GREENLIST_USER_IDS}",SUNO_TOKEN="${SUNO_TOKEN}",WEBHOOK_BOT="${WEBHOOK_BOT}",WEBHOOK_SEND_TO="${WEBHOOK_SEND_TO}",STRATASLIMS_LOG_LEVEL="INFO"

; Healthcheck via preflight command before declaring process started
; Use eventlistener or startsecs to allow network to become ready; our main.py also waits.

; Logging with rotation managed by supervisord
stdout_logfile=/var/log/supervisor/strataslims.log
stderr_logfile=/var/log/supervisor/strataslims_err.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile_backups=5

; Keep environment and locale predictable
environment=PYTHONUNBUFFERED="1",LC_ALL="C.UTF-8",LANG="C.UTF-8"
